# etc/sv/openvpn/run
#!/bin/sh

[ -r conf ] && . ./conf
exec /usr/bin/openvpn --config /etc/openvpn/client/client.ovpn $OPTS
